# toggle_theme.service
[Unit]
Description=Toggle theme

[Service]
ExecStart=/home/cuikho210/Documents/scripts/toggle-theme/toggle_theme.bash
